# deploy/animaworks.service
# AnimaWorks systemd user/service unit template.
# Edit the paths below to match your environment before installing.
# Typical install: copy to ~/.config/systemd/user/animaworks.service and run systemctl --user daemon-reload.

[Unit]
Description=AnimaWorks - Digital Anima Framework
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=20

[Service]
Type=simple
WorkingDirectory=/path/to/animaworks
ExecStart=/path/to/animaworks/.venv/bin/animaworks start --foreground
ExecStop=/path/to/animaworks/.venv/bin/animaworks stop
Restart=on-failure
RestartSec=10
RestartPreventExitStatus=3
LimitNOFILE=524288
Environment=HOME=/path/to/home
Environment=PATH=/path/to/home/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TOKENIZERS_PARALLELISM=false

[Install]
WantedBy=default.target
